# crontab
0 8 * * * Rscript /app.R >> /var/log/app.log 2>&1
