# crontab
0 * * * * /usr/local/bin/Rscript /app/app.R >> /var/log/app.log 2>&1
